# rootfs/etc/s6-overlay/s6-rc.d/zegy/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Zegy Sensor Manager..."

export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"

exec node /app/backend/dist/index.js
